# deploy/systemd/caddy-proxy.service
[Unit]
Description=Caddy shared reverse proxy
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/caddy-proxy/deploy.sh
